# DSmigrate Files/Library/LaunchDaemons/com.318.dsmigrate.launchd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Disabled</key>
	<false/>
	<key>HopefullyExitsLast</key>
	<true/>
	<key>Label</key>
	<string>com.318.dsmigrate</string>
	<key>OnDemand</key>
	<false/>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>-c</string>
		<string>/usr/local/sbin/dsmigrate -u;
/bin/mv /Library/LaunchDaemons/com.318.dsmigrate.launchd.plist /private/tmp/com.318.dsmigrate.launchd.plist;
/bin/rm /usr/local/sbin/dsmigrate;
/bin/rm /usr/local/sbin/dsmigrate.conf;
/bin/launchctl unload /private/tmp/com.318.dsmigrate.launchd.plist;
sleep 60;
exit 0</string>
	</array>
</dict>
</plist>
